# private_dot_config/systemd/user/wob.socket
[Unit]
Description=wob — Wayland Overlay Bar
Documentation=man:wob(1)

[Socket]
ListenFIFO=%t/wobpipe

[Install]
WantedBy=sway-session.target
