# private_dot_config/systemd/user/redshift-gtk.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/redshift-gtk -m wayland

[Install]
WantedBy=
WantedBy=sway-session.target
